# deb/xenial/debian/postinst
#!/bin/bash
set -e

# Setting package version
VERSION='%VERSION%'
RUN_VERSION=$2
ARG=$1

/usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 777 /tmp' > /dev/null 2>&1
/usr/sbin/chroot /opt/gluu-server /bin/su - root -c '/tmp/system_user.list' > /dev/null 2>&1
/usr/sbin/chroot /opt/gluu-server /bin/su - root -c '/tmp/system_group.list' > /dev/null 2>&1

/bin/cp -a /etc/resolv.conf /opt/gluu-server/etc/ > /dev/null 2>&1

if [ -e /etc/init.d/gluu-server ]; then
    update-rc.d gluu-server defaults
fi

if [ "$1" = "configure" ]; then
  if [ -z "$RUN_VERSION" ] || [ "$RUN_VERSION" = "$VERSION" ]; then
    echo "Starting gluu-server ..."
    sleep 1
    /etc/init.d/gluu-server start
    sleep 4
    /etc/init.d/gluu-server login && /opt/gluu/bin/install.py && exit 0
  fi
fi
 
exit 0
