# deb/xenial/debian/postinst
#!/bin/bash
set -e

# Setting package version
VERSION='%VERSION%'
RUN_VERSION=$2
ARG=$1

/usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 777 /tmp' > /dev/null 2>&1
/usr/sbin/chroot /opt/gluu-server /bin/su - root -c '/tmp/system_user.list' > /dev/null 2>&1
/usr/sbin/chroot /opt/gluu-server /bin/su - root -c '/tmp/system_group.list' > /dev/null 2>&1

/bin/cp -a /etc/resolv.conf /opt/gluu-server/etc/ > /dev/null 2>&1

if [ -e /etc/init.d/gluu-server ]; then
    update-rc.d gluu-server defaults
fi

if [ "$1" = "configure" ]; then
  if [ -z "$RUN_VERSION" ] || [ "$RUN_VERSION" = "$VERSION" ]; then
    echo "Starting gluu-server ..."
    sleep 1
    /etc/init.d/gluu-server start
    sleep 5
    /usr/sbin/chroot /opt/gluu-server /bin/su - root -c '/opt/gluu/bin/install.py'
    rm -rf /opt/gluu-server/etc/hosts
    rm -rf /opt/gluu-server/etc/hostname
    mv /opt/gluu-server/etc/hosts.original /opt/gluu-server/etc/hosts
    mv /opt/gluu-server/etc/hostname.original /opt/gluu-server/etc/hostname
  else
    if [ -n "$RUN_VERSION" ] && [ "$RUN_VERSION" != "$VERSION" ]; then
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'useradd ldap -d /home/ldap'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'useradd jetty -d /home/jetty'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'useradd node -d /home/node'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'useradd radius -d /opt/gluu/radius'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'groupadd gluu'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'gpasswd -a ldap gluu > /dev/null'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'gpasswd -a node gluu > /dev/null'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'gpasswd -a jetty gluu > /dev/null'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'gpasswd -a ldap adm > /dev/null'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'gpasswd -a radius gluu > /dev/null'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'gpasswd -a ldap adm > /dev/null'         
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R jetty:jetty /opt/jetty*'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R jetty:jetty /opt/gluu/jetty'      
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R jetty:jetty /opt/shibboleth-idp'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R root:gluu /etc/certs'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R jetty:jetty /etc/certs/oxauth-keys.jks'      
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R ldap:ldap /opt/opendj'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R node:node /opt/node*'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R node:node /opt/gluu/node'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown -R radius:ldap /opt/gluu/radius'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown ldap:ldap /home/ldap'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown node:node /home/node'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /home/jetty'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 644 /etc/profile'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod ga+w /tmp'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:gluu /var/gluu/photos'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:gluu /var/gluu/identity/removed'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chwon root:gluu /var/gluu/identity/cr-snapshots'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 775 /var/run/jetty'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chgrp jetty /var/run/jetty'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/httpd.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/httpd.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/httpd.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/httpd.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/shibIDP.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/shibIDP.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/shibIDP.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/shibIDP.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/idp-encryption.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/idp-encryption.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/idp-encryption.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/idp-encryption.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/idp-signing.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/idp-signing.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/idp-signing.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/idp-signing.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown ldap:ldap /etc/certs/passport-sp.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/passport-sp.key.orig'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown ldap:ldap /etc/certs/passport-sp.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/passport-sp.key'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/shibIDP.pkcs12'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/shibIDP.pkcs12'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/shibIDP.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 700 /etc/certs/shibIDP.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /install/community-edition-setup/output/oxauth-keys.json'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 500 /install/community-edition-setup/output/oxauth-keys.json'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 644 /etc/hostname'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 644 /etc/hosts'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown ldap:ldap /home/ldap/.pw'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:gluu /etc/certs'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:gluu /etc/gluu/conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:gluu /var/gluu'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 440 /etc/certs'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod a+x /etc/certs'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown radius:gluu /etc/certs/gluu-radius.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown radius:gluu /etc/certs/gluu-radius.private-key.pem'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown jetty:jetty /etc/certs/oxauth-keys.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 660 /etc/certs/oxauth-keys.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod +x /etc/init.d/opendj'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /etc/default/oxauth'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod +x /etc/init.d/oxauth'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /usr/lib/tmpfiles.d/jetty.conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 644 /usr/lib/tmpfiles.d/jetty.conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /etc/default/identity'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod +x /etc/init.d/identity'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /usr/lib/tmpfiles.d/jetty.conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 644 /usr/lib/tmpfiles.d/jetty.conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /etc/default/idp'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod +x /etc/init.d/idp'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /usr/lib/tmpfiles.d/jetty.conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 644 /usr/lib/tmpfiles.d/jetty.conf'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:root /etc/default/passport'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod +x /etc/init.d/passport'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod +x /etc/init.d/gluu-radius'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown root:gluu /etc/gluu/conf/radius/'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown radius:gluu /etc/certs/gluu-radius.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown radius:gluu /etc/certs/gluu-radius.private-key.pem'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chown radius:gluu /etc/certs/oxauth-keys.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 660 /etc/certs/oxauth-keys.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 660 /etc/certs/gluu-radius.jks'
      /usr/sbin/chroot /opt/gluu-server /bin/su - root -c 'chmod 660 /etc/certs/gluu-radius.private-key.pem'
      /etc/init.d/gluu-server stop
      /etc/init.d/gluu-server start 
    fi
  fi
fi

exit 0
